# settings/nursing_export/centos/crontab
5 1 * * * /usr/local/ienursing2_cron/export_nursing.sh